# root/etc/s6-overlay/s6-rc.d/init-config-mosquitto/run
#!/usr/bin/with-contenv bash

[[ ! -f /config/mosquitto.conf ]] &&
	cp /defaults/mosquitto.conf /config/mosquitto.conf

# permissions
chown -R abc:abc \
	/config/
